# services/boot-diagnostics.service
[Unit]
Description=Skipper Boot Diagnostics Service
After=multi-user.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /home/pi/work/skipper/services/boot_diagnostics.py
WorkingDirectory=/home/pi/work/skipper
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
